# recipes-core/edgeos-identity/files/edgeos-uuid-generate.service
[Unit]
Description=EdgeOS Device UUID Generation
DefaultDependencies=no
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/generate-uuid.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=sysinit.target
