# recipes-support/nvidia-container-config/files/edgeos-cdi-generate.service
[Unit]
Description=EdgeOS CDI Spec Generation for NVIDIA GPUs
After=local-fs.target systemd-udev-settle.service
Before=containerd.service docker.service
# Only run if nvidia-cdi-refresh.service doesn't exist (fallback)
ConditionPathExists=!/lib/systemd/system/nvidia-cdi-refresh.service
# Ensure l4t.csv exists before running
ConditionPathExists=/etc/nvidia-container-runtime/host-files-for-container.d/l4t.csv

[Service]
Type=oneshot
# Wait a bit for device nodes to appear
ExecStartPre=/bin/sleep 2
# Create directory if needed
ExecStartPre=/bin/mkdir -p /etc/cdi
# Generate CDI spec using CSV mode with explicit CSV file list
# Includes: devices.csv (meta-tegra), drivers.csv (meta-tegra), l4t.csv (ours),
#           devices-wendyos.csv (sysfs entries), and l4t-deepstream.csv (if present)
ExecStart=/bin/sh -c '\
    ARGS="--mode=csv --output=/etc/cdi/nvidia.yaml"; \
    ARGS="$ARGS --csv.file=/etc/nvidia-container-runtime/host-files-for-container.d/devices.csv"; \
    ARGS="$ARGS --csv.file=/etc/nvidia-container-runtime/host-files-for-container.d/drivers.csv"; \
    ARGS="$ARGS --csv.file=/etc/nvidia-container-runtime/host-files-for-container.d/l4t.csv"; \
    ARGS="$ARGS --csv.file=/etc/nvidia-container-runtime/host-files-for-container.d/devices-wendyos.csv"; \
    if [ -f /etc/nvidia-container-runtime/host-files-for-container.d/l4t-deepstream.csv ]; then \
        ARGS="$ARGS --csv.file=/etc/nvidia-container-runtime/host-files-for-container.d/l4t-deepstream.csv"; \
    fi; \
    /usr/bin/nvidia-ctk cdi generate $ARGS'
# Post-process CDI spec to fix GStreamer plugin paths for container compatibility
# (Yocto uses /usr/lib/ but containers expect /usr/lib/aarch64-linux-gnu/)
ExecStartPost=/usr/bin/fix-cdi-gstreamer-paths.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
