# recipes-network/gadget-network-config/files/usb0-route.service

[Unit]
Description=Set default route via USB host (from dnsmasq leases)
After=network-online.target gadget-dnsmasq.service
Wants=network-online.target

[Service]
Type=oneshot
Environment=LEASE_FILE=/var/lib/misc/dnsmasq.leases
Environment=DEV=usb0
ExecStart=/usr/bin/usb0-route-from-leases.sh

