# recipes-extended/edgeos-user/files/edgeos-user-setup.service
[Unit]
Description=EdgeOS User Home Directory Setup (First Boot)
Documentation=man:systemd.service(5)
# Must run after /data is mounted and /home is bind mounted
RequiresMountsFor=/data /home
After=data.mount home.mount
# Run before user sessions start
Before=systemd-user-sessions.service
# Only run once (script checks for flag file)
ConditionPathExists=!/data/.edgeos-user-setup-done

[Service]
Type=oneshot
ExecStart=/usr/sbin/edgeos-user-setup.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
