# recipes-core/systemd-mount-containerd/files/var-lib-containerd.mount
[Unit]
Description=Bind mount /var/lib/containerd from persistent data partition
Documentation=man:systemd.mount(5)

# Ensure data partition is mounted and filesystem is fully expanded
RequiresMountsFor=/data
After=data.mount

# Wait for Mender data partition resize to complete (if enabled)
# This is safe even if mender-systemd-growfs-data.service doesn't exist
After=mender-systemd-growfs-data.service

# Must be mounted before containerd starts
Before=containerd.service
DefaultDependencies=no

[Mount]
What=/data/containerd
Where=/var/lib/containerd
Type=none
Options=bind,nosuid,nodev,x-systemd.mkdir

[Install]
WantedBy=local-fs.target
