# recipes-core/wendyos-identity/files/wendyos-uuid-generate.service
[Unit]
Description=WendyOS Device UUID Generation
DefaultDependencies=no
After=local-fs.target data.mount wendyos-etc-binds.service
Wants=wendyos-etc-binds.service

[Service]
Type=oneshot
ExecStart=/usr/bin/generate-uuid.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=sysinit.target
