# recipes-devtools/wendyos-containerd-registry/files/wendyos-dev-registry.service
[Unit]
Description=WendyOS Development Container Registry
Documentation=https://github.com/mihai-chiorean/containerd-registry
After=containerd.service network-online.target
Requires=containerd.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=REGISTRY_NAMESPACE=default
Environment=REGISTRY_NAME=wendyos-dev-registry
Environment=REGISTRY_IMAGE=wendyos/containerd-registry:v1.0.7
Environment=LISTEN_ADDRESS=0.0.0.0:5000

# Pre-flight checks: Verify image exists, import if missing (self-healing)
ExecStartPre=/bin/sh -c '\
  if ! /usr/bin/ctr -n ${REGISTRY_NAMESPACE} images ls | /bin/grep -q ${REGISTRY_IMAGE}; then \
    echo "Registry image not found, importing from archive..."; \
    /usr/bin/ctr -n ${REGISTRY_NAMESPACE} images import /usr/share/wendyos/offline-images/containerd-registry-latest.tar; \
    /usr/bin/ctr -n ${REGISTRY_NAMESPACE} images tag wendyos/containerd-registry:v1.0.7 wendyos/containerd-registry:latest || true; \
    echo "Registry image imported successfully"; \
  fi'
# Cleanup any existing container
ExecStartPre=-/usr/bin/nerdctl rm -f ${REGISTRY_NAME}

# Start the registry container with nerdctl for journald logging
ExecStart=/usr/bin/nerdctl run -d \
    --name=${REGISTRY_NAME} \
    --namespace=${REGISTRY_NAMESPACE} \
    --network=host \
    --mount type=bind,src=/run/containerd/containerd.sock,dst=/run/containerd/containerd.sock \
    --env LISTEN_ADDRESS=${LISTEN_ADDRESS} \
    --env CONTAINERD_NAMESPACE=${REGISTRY_NAMESPACE} \
    --env LOG_FORMAT=json \
    --log-driver=journald \
    --log-opt=tag=${REGISTRY_NAME} \
    --restart=no \
    ${REGISTRY_IMAGE}

# Graceful shutdown
ExecStop=/usr/bin/nerdctl stop -t 30 ${REGISTRY_NAME}
ExecStopPost=-/usr/bin/nerdctl rm -f ${REGISTRY_NAME}

# No automatic restart for oneshot services - handled by wendy-agent if needed
TimeoutStartSec=60s
TimeoutStopSec=45s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wendyos-dev-registry

# Security
NoNewPrivileges=true

[Install]
# Disabled by default - started on-demand by wendy-agent
# To enable auto-start on boot: systemctl enable wendyos-dev-registry
WantedBy=multi-user.target
